# robot/crontab
0 2,5,8,11,14,17,20,23 * * * python /root/silver_weather/SilverWeahter/robot/jarvis.py >/dev/null 2>&1 &